# infra/systemd/trendscope-postgres-backup.timer
[Unit]
Description=Run The Trend Scope PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
